# devops/systemd/welllabs-frontend.service
[Unit]
Description=WellLabs Frontend — React SPA via serve (port 3000)
Documentation=https://github.com/vercel/serve
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/welllabs/frontend/current

# Serve React build in SPA mode (-s = single-page app, serves index.html for all 404s)
ExecStart=/usr/local/bin/serve -s /opt/welllabs/frontend/current -l 127.0.0.1:3000

# Environment
Environment=NODE_ENV=production

Restart=always
RestartSec=5

StandardOutput=append:/opt/welllabs/logs/frontend.log
StandardError=append:/opt/welllabs/logs/frontend-error.log

[Install]
WantedBy=multi-user.target