# devops/systemd/welllabs-frontend.service
[Unit]
Description=WellLabs Frontend — React SPA via serve (port 3000)
Documentation=https://github.com/vercel/serve
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/welllabs/frontend/current
ExecStart=/usr/bin/serve -s /opt/welllabs/frontend/current -l 3000
Restart=always
RestartSec=5
StandardOutput=append:/opt/welllabs/logs/frontend.log
StandardError=append:/opt/welllabs/logs/frontend-error.log

[Install]
WantedBy=multi-user.target